# etc/supervisor/conf-available.d/monitor.conf.in
[group:monitor]
programs=x11,vncserver,chromium

[program:x11]
command=@@INSTALL_DIR@@/bin/x11.bash
environment=STARTUP="@@INSTALL_DIR@@/bin/ui.bash"
autorestart=true
stopwaitsecs=10
startretries=9223372036854775807
priority=2

redirect_stderr=true
stdout_logfile=@@MONITOR_LOG_DIR@@/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=NONE

[program:vncserver]
command=@@INSTALL_DIR@@/bin/vncserver.bash
autorestart=true
stopwaitsecs=10
startretries=9223372036854775807
priority=2

redirect_stderr=true
stdout_logfile=@@MONITOR_LOG_DIR@@/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=NONE

[program:chromium]
command=@@INSTALL_DIR@@/bin/chromium.bash
autorestart=true
stopwaitsecs=10
startretries=9223372036854775807
priority=3

redirect_stderr=true
stdout_logfile=@@MONITOR_LOG_DIR@@/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=NONE

